# system_files/clock.timer
[Unit]
Description=Start Clock

[Timer]
OnBootSec=1

[Install]
WantedBy=timers.target
