# system_files/clock.service
[Unit]
Description=Start Clock

[Install]
WantedBy=clock.timer

[Service]
Type=simple
ExecStart=/usr/local/sbin/clock
